# deb/postinst
#!/bin/sh

mkdir -p /var/log/paparazzi

systemctl daemon-reload
systemctl enable paparazzi.service
systemctl start paparazzi.service
