# deb/prerm
#!/bin/sh

systemctl stop paparazzi.service
systemctl disable paparazzi.service
